# print-server/print-server.service
[Unit]
Description=Penkey POS Print Server
After=network.target

[Service]
Type=simple
User=jimmy
WorkingDirectory=/home/jimmy/print-server
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/jimmy/print-server/venv/bin/python /home/jimmy/print-server/print-server/print_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=penkey-print-server

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp /home/jimmy/print-server /dev /home/jimmy/print-server/.git /home/jimmy/print-server/print-server
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
